# auto_down.cron
*/1 * * * /usr/src/OneGuyApp/auto_down.sh
